# videok/autostart/video.desktop
[Desktop Entry]
Name=Play a video
Exec=mpv --fs