# services/shutdown.service
[Unit]
Description=Shutdown RPi at button push
After=Network.target

[Service]
ExecStart=/bin/shutdown

[Install]
WantedBy=multi-user.target